# src/skel/etc/supervisor.d/program-collectd.conf
[program:collectd]
command = /usr/local/bin/launch-collectd.sh
