# src/config/program-crond.conf
[program:crond]
command = /usr/sbin/crond -m off -s -n
